# crontab
0 0 * * * python3 /Users/kevinkim/Desktop/untitled folder/flask/fetch.py
